# chargen.service
[Unit]
Description=Character Generator
After=network.target

[Service]
ExecStart=/usr/local/bin/chargen

[Install]
WantedBy=multi-user.target
